# etc/uci-defaults/tch_9001-upgrade_sts
#!/bin/sh
uci set network.fonopen.ipv6=1
uci commit network
uci -q get mobiled.device_defaults.minimum_attach_delay && exit 0
uci set mobiled.device_defaults.minimum_attach_delay=0
uci set mobiled.device_defaults.maximum_attach_delay=10
uci set mobiled.device_defaults.detach_mode=none
uci set mobiled.@device[0].minimum_attach_delay=0
uci set mobiled.@device[0].maximum_attach_delay=10
uci set mobiled.@device[0].detach_mode=none
uci delete mobiled.@mobiled_state[11].entryexits
uci set mobiled.@mobiled_state[12].entryexits=DataSessionSetup-EntryExit
uci set mobiled.@mobiled_state[10].timeout=60
uci commit mobiled
uci set wansensing.global.autofailovermaxwait=10
uci commit wansensing
uci set tod.global.tod_enabled=1
uci set tod.global.tracelevel=4
uci set tod.global.time_change_monfreq=10
uci commit tod
uci -q get mobiled.globals.pdn_retry_timer && exit 0
uci set mobiled.globals.pdn_retry_timer=720
uci commit mobiled
uci add_list parental.general.lan_intf='Guest1'
uci add_list parental.general.lan_intf='Guest1-5GHz'
uci commit parental
exit 0
