# etc/init.d/geolocation
#!/bin/sh /etc/rc.common
START=99
USE_PROCD=1

start_service() {
    config_load geolocation
    config_get_bool enabled global enabled 0
    [ $enabled -eq 0 ] && exit
    /usr/bin/geolocation
}
