# usr/lib/opkg/info/geolocation-tch.list
/usr/bin/geolocation
/etc/init.d/geolocation
